# recipes-newbs/newbs-nvram/files/newbs-nvram.service.in
[Unit]
Description=NEWBS NVRAM Bind mounts
DefaultDependencies=false
RequiresMountsFor=@NEWBS_NVRAM_DIR@
Before=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
StandardOutput=journal+console
ExecStart=/usr/lib/newbs/newbs-nvram

[Install]
WantedBy=local-fs.target
